# contrib/histkit-scan.timer
[Unit]
Description=Run histkit scan periodically

[Timer]
OnBootSec=5m
OnUnitActiveSec=12h
Persistent=true

[Install]
WantedBy=timers.target
