# contrib/histkit-scan.service
[Unit]
Description=Scan and index shell history with histkit

[Service]
Type=oneshot
ExecStart=%h/.local/bin/histkit scan --config %h/.config/histkit/config.toml
